# Splunk-Detection-Lab/detections/11_registry_run_key.spl
# Detection: Registry Run Key Persistence
# MITRE ATT&CK: T1547.001 - Boot or Logon Autostart Execution: Registry Run Keys
# Severity: High
# Data Source: Sysmon (Event ID 13 - Registry value set)

index=sysmon EventCode=13
  (TargetObject="*\\CurrentVersion\\Run*" OR TargetObject="*\\CurrentVersion\\RunOnce*")
| table _time, host, Image, TargetObject, Details
| eval severity="high"
| eval mitre_technique="T1547.001"
